# scripts/com.dudu.daily-price-upload.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.dudu.daily-price-upload</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>/Users/dudu/Desktop/trae/重点日常项目/【GTM】新机售价监控分析系统/scripts/upload-daily-price-excel.sh</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>10</integer>
    <key>Minute</key>
    <integer>20</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>/tmp/daily-price-upload.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/daily-price-upload.err</string>
</dict>
</plist>
